# service/autoscroll-indicator.service
[Unit]
Description=autoscroll-linux — Wayland overlay indicator
Documentation=man:autoscroll(1)
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/autoscroll-indicator
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
