# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Create dedicated service user if it doesn't exist
        if ! getent passwd autoscroll >/dev/null 2>&1; then
            adduser --system --no-create-home --group autoscroll 2>/dev/null || true
        fi

        # Reload udev rules so /dev/uinput permissions apply
        udevadm control --reload-rules 2>/dev/null || true
        udevadm trigger --subsystem-match=input 2>/dev/null || true

        # Enable and start the service
        systemctl --no-reload enable autoscroll.service 2>/dev/null || true
        systemctl restart autoscroll.service 2>/dev/null || true
        ;;
esac
